# com.rich.torrentremover.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.rich.removeTorrents</string>

    <key>ProgramArguments</key>
    <array>
      <string>~/bin/removeTorrents</string>
    </array>

    <!-- Run once when the agent is loaded -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- KeepAlive can be configured, but we leave it off for a simple run -->
  </dict>
</plist>
